# server/pirobot.service
[Unit]
Description=PiRobot Daemon

[Service]
Restart=always
Type=simple
User=www-data
ExecStart=/usr/local/bin/pirobotd -c picaterpillar

[Install]
WantedBy=multi-user.target
